# systemd/collab-cluster-tracker.service
[Unit]
Description=BitTorrent v2 prototype — tracker (announce + catalog), port 6969

[Service]
WorkingDirectory=%h/collab-cluster-experiment
ExecStart=python3 bittorrent_tracker.py
Restart=on-failure

[Install]
WantedBy=default.target
